# actualbudget2mqtt/rootfs/etc/services.d/actualbudget2mqtt/finish
#!/command/with-contenv bashio

echo "[FINISH] ActualBudget2MQTT finished."