# actualbudget2mqtt/rootfs/etc/services.d/actualbudget2mqtt/finish
#!/bin/sh
bashio::log.info "Stopping ActualBudget2MQTT service ..."
