# actualbudget2mqtt/rootfs/etc/services.d/actualbudget2mqtt/finish
#!/usr/bin/env bash
echo "[FINISH] ActualBudget2MQTT finished."