# actualbudget2mqtt/rootfs/etc/services.d/actualbudget2mqtt/run

#!/command/with-contenv bashio
set -euo pipefail

bashio::log.info "Starting ActualBudget2MQTT service ..."
exec node /opt/actualbudget2mqtt/index.js
