# actualbudget2mqtt/rootfs/etc/services.d/actualbudget2mqtt/run
#!/bin/sh
set -eu

bashio::log.info "Starting ActualBudget2MQTT service ..."
exec node /opt/actualbudget2mqtt/index.js
