# actualbudget2mqtt/rootfs/etc/services.d/actualbudget2mqtt/run
#!/command/with-contenv bashio

echo "[RUN] Starting ActualBudget2MQTT..."
node /actualbudget2mqtt.js