# actualbudget2mqtt/rootfs/etc/services.d/actualbudget2mqtt/run
#!/usr/bin/env bash
echo "[RUN] Starting ActualBudget2MQTT..."
node /actualbudget2mqtt.js