# technitium_dns/rootfs/etc/services.d/technitium/finish
#!/usr/bin/with-contenv bash
echo "[Technitium DNS] Dienst beendet."
