# technitium_dns/rootfs/etc/services.d/technitium_dns/finish
#!/usr/bin/with-contenv sh
echo "[WARN] Technitium DNS wurde beendet. Container wird ggf. neu gestartet..."