# technitium_dns/rootfs/etc/services.d/technitium_dns/finish
#!/usr/bin/with-contenv sh
. /opt/scripts/common.sh
log_warn "Technitium DNS wurde beendet. Container wird ggf. neu gestartet..."
