# technitium_dns/rootfs/etc/services.d/technitium_dns/finish
#!/usr/bin/with-contenv bash
source /opt/scripts/common.sh
log_warn "Technitium DNS wurde beendet. Container wird neu gestartet..."
