# technitium_dns/rootfs/etc/services.d/technitium_dns/run
#!/usr/bin/with-contenv bash
set -e

/rootfs/etc/services.d/technitium_dns/scripts/log.sh "[INFO] Starte Technitium DNS Add-on"

/rootfs/etc/services.d/technitium_dns/scripts/init.sh

/rootfs/etc/services.d/technitium_dns/scripts/healthcheck.sh

/rootfs/etc/services.d/technitium_dns/scripts/log.sh "[INFO] Starte Technitium DNS Server"
exec /app/TechnitiumDnsServer
