# technitium_dns/rootfs/etc/services.d/technitium_dns/run
#!/usr/bin/with-contenv sh
. /opt/scripts/common.sh
. /opt/scripts/init.sh

log_info "Starte Technitium DNS Server..."
log_info "Persistenter Ordner: $APPDATAPATH"

# Technitium DNS starten, nutzt APPDATAPATH automatisch
exec /app/dns/dnsserver