# technitium_dns/rootfs/etc/services.d/technitium_dns/run
#!/usr/bin/with-contenv sh
# Starte Technitium DNS

. /opt/scripts/common.sh
. /opt/scripts/init.sh

log_info "Starte Technitium DNS Server..."

# Standardstart aus dem offiziellen Image
exec /app/dns/dnsserver --appdatafolder /data