# technitium_dns/rootfs/etc/services.d/technitium_dns/run
#!/usr/bin/with-contenv sh
. /opt/scripts/common.sh
. /opt/scripts/init.sh

log_info "Starte Technitium DNS Server (Custom Multi-Arch Image)..."
log_info "Persistenter Ordner: /data"

exec /app/dns/dnsserver --appdatafolder /data
